# systemd/myopsida.service
[Unit]
Description=Update Cloudflare DNS records with current public IP
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/myopsida/config.env
ExecStart=/usr/local/bin/myopsida \
  -zone-id ${CF_ZONE_ID} \
  -record-id-v4 ${CF_RECORD_ID_V4} \
  -record-id-v6 ${CF_RECORD_ID_V6} \
  -name ${CF_RECORD_NAME}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=myopsida

[Install]
WantedBy=multi-user.target
